# codex_bridge_app/rootfs/etc/s6-overlay/s6-rc.d/codex-bridge-init/up
/usr/local/libexec/codex-bridge/initialize.sh
